# upstart/simplelogin.conf
pre-start script
end script

respawn
respawn limit 15 5

start on ((filesystem
           and started dbus)
          or runlevel PREVLEVEL=S)


stop on runlevel [016]

emits login-session-start
emits desktop-session-start
emits desktop-shutdown

script
/usr/bin/simplelogin --user phablet --exec /usr/bin/kwinwrapper
end script
